# trading_tool/deploy/preview-watchdog.conf
[program:preview-watchdog]
command=bash /workspace/trading_tool/preview_watchdog.sh
directory=/workspace/trading_tool
autostart=true
autorestart=true
startsecs=5
startretries=10
stopsignal=INT
stopwaitsecs=10
stdout_logfile=/tmp/preview-watchdog.log
stderr_logfile=/tmp/preview-watchdog.log
redirect_stderr=true
environment=PATH="%(ENV_PATH)s",NODE_OPTIONS=""
