# trading_tool/deploy/preview-8000.conf
[program:preview-8000]
command=python3.11 /workspace/trading_tool/web_app.py
directory=/workspace/trading_tool
autostart=true
autorestart=true
startsecs=5
startretries=10
stopsignal=INT
stopwaitsecs=10
stopasgroup=true
killasgroup=true
stdout_logfile=/tmp/preview-8000.log
stderr_logfile=/tmp/preview-8000.log
redirect_stderr=true
environment=PATH="%(ENV_PATH)s",NODE_OPTIONS=""
